# apps/kinton/supervisor/mosquitto.conf
[program:mosquitto]
command=/kinton/scripts/start-mosquitto.sh
autostart=true
autorestart=true
